# renfield-satellite/systemd/renfield-satellite.service
[Unit]
Description=Renfield Satellite Voice Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=evdb
Group=evdb
WorkingDirectory=/opt/renfield-satellite
ExecStart=/opt/renfield-satellite/venv/bin/python -m renfield_satellite config/satellite.yaml
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Allow access to hardware
SupplementaryGroups=spi gpio i2c

# Resource limits (Pi Zero 2 W has 512MB RAM)
MemoryMax=300M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=renfield-satellite

[Install]
WantedBy=multi-user.target
